# automation/systemd/local-data-stack.service
[Unit]
Description=Local Data Stack Analytics Pipeline
Documentation=https://github.com/flucido/local-data-stack
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=datapipeline
Group=datapipeline
WorkingDirectory=/Users/flucido/projects/local-data-stack

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="LOG_LEVEL=INFO"
Environment="PATH=/Users/flucido/projects/local-data-stack/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Main execution command
ExecStart=/Users/flucido/projects/local-data-stack/.venv/bin/python3 /Users/flucido/projects/local-data-stack/scripts/run_pipeline.py --stage all

# Pre-execution validation
ExecStartPre=/bin/sh -c 'test -f /Users/flucido/projects/local-data-stack/.env || (echo ".env file missing" && exit 1)'
ExecStartPre=/bin/sh -c 'test -f /Users/flucido/projects/local-data-stack/oss_framework/data/oea.duckdb || echo "Creating new DuckDB database"'

# Restart strategy
Restart=on-failure
RestartSec=30s
StartLimitBurst=3
StartLimitIntervalSec=600

# Resource limits (based on 2026 best practices)
MemoryMax=4G
MemoryHigh=3G
CPUQuota=200%
TasksMax=512

# I/O limits (prevent disk thrashing)
IOWeight=500
IOReadBandwidthMax=/dev/sda 50M
IOWriteBandwidthMax=/dev/sda 50M

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/Users/flucido/projects/local-data-stack/oss_framework/data
ReadWritePaths=/Users/flucido/projects/local-data-stack/oss_framework/logs
ReadWritePaths=/Users/flucido/projects/local-data-stack/rill_project/data
ReadWritePaths=/tmp/duckdb_temp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=local-data-stack

[Install]
WantedBy=multi-user.target
